# tarsnap.service
[Unit]
Description=Tarsnap backup

[Service]
ExecStart=/usr/local/sbin/backup
